# debian-recursor/pdns.conffiles
/etc/powerdns/recursor.conf
/etc/init.d/pdns-recursor
