# debian-pdns/pdns.prerm
#!/bin/sh

set -e

if [ "$1" = "remove" ] ; then
	if command -v invoke-rc.d >/dev/null 2>&1; then
		invoke-rc.d pdns stop
	else
		/etc/init.d/pdns stop
	fi
fi
