# build-scripts/debian-authoritative-semistatic/pdns.conffiles
/etc/powerdns/pdns.conf
/etc/init.d/pdns
