# build-scripts/debian-authoritative-semistatic/pdns.prerm
#!/bin/sh
# prerm script for pdns

set -e

case "$1" in
    remove|upgrade|deconfigure)
        if command -v invoke-rc.d >/dev/null 2>&1; then
            invoke-rc.d pdns stop
        else
            /etc/init.d/pdns stop
        fi
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
