# debian/pdns-backend-mysql.prerm
#!/bin/sh
#
# Pre removal

set -e

# Stop pdns.
if [ -x "/etc/init.d/pdns" ]; then
	invoke-rc.d pdns stop || exit $?
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

