# debian/postrm
#!/bin/sh

set -e

update_init() {
	update-rc.d pdns remove > /dev/null
}


if [ "$1" = "remove" ] ; then
	ldconfig
elif [ "$1" = "purge" ] ; then
	update_init
fi

