# pdns/debian/postinst
#!/bin/sh

set -e

if [ "$1" != "configure" ] ; then
	exit 0
fi

update-rc.d pdns defaults > /dev/null

if [ -n "$2" ] ; then
	/etc/init.d/pdns restart
else
	/etc/init.d/pdns start
fi

exit 0

