# pdns/debian/postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ] ; then
	update-rc.d pdns remove > /dev/null
	userdel pdns
fi

# TODO: check if we need to remove files from /var/run

exit 0

