# builder-support/debian/dnsdist/ubuntu-trusty/dnsdist.upstart
description "dnsdist - A powerful DNS loadbalancer"

start on runlevel [2345]
stop on runlevel [!2345]

author "Pieter Lexis <pieter.lexis@powerdns.com>"

# Keep the --supervised option when modifying this
exec /usr/bin/dnsdist --supervised
